# leadaudit.service
[Unit]
Description=LeadAudit Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/lead-audit
ExecStart=/root/lead-audit/venv/bin/python /root/lead-audit/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
